# setup/etc/systemd/system/storage.service
[Unit]
Description=Storage Monitor Service

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/local/bin/storage.sh 2>&1'

