# setup/etc/systemd/system/storage.path
[Unit]
Description=Storage Module Monitor

[Path]
PathChanged=/usr/share/issabel/storage/

[Install]
WantedBy=multi-user.target
